# deploy/systemd/banking-backend-celery-beat.service
[Unit]
Description=SafaPay Bank - Celery beat scheduler
After=network.target

[Service]
Type=simple

# Adjust these for your server:
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/banking-backend
EnvironmentFile=/opt/banking-backend/.env

ExecStart=/opt/banking-backend/venv/bin/celery -A config beat -l INFO --scheduler django_celery_beat.schedulers:DatabaseScheduler

Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target

